# proxy.service
[Unit]
Description=SNI Proxy (podman container)
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/home/admin/py_proxy /home/admin/ssl

[Service]
Type=simple
Restart=always
RestartSec=5

ExecStartPre=-/usr/bin/podman rm -f py-proxy

ExecStart=/usr/bin/podman run \
    --name py-proxy \
    --rm \
    --network=host \
    -v /home/admin/py_proxy/config.json:/app/config.json:Z,ro \
    -v /home/admin/ssl:/home/admin/ssl:Z,ro \
    localhost/py-proxy:latest

ExecStop=/usr/bin/podman stop py-proxy
ExecReload=/usr/bin/podman kill --signal HUP py-proxy

StandardOutput=journal
StandardError=journal
SyslogIdentifier=py-proxy

[Install]
WantedBy=multi-user.target
